# systemd/tpm2-ssh-agent.service
[Unit]
Description=TPM2 SSH agent

[Service]
ExecStart=/usr/bin/tpm2-ssh-agent --key /PATH/TO/KEY -f
Type=simple

[Install]
WantedBy=default.target
