# debian/tpm2-ssh-agent.socket
[Unit]
Description=TPM2 SSH agent socket

[Socket]
ListenStream=/run/tpm2_ssh_agent.socket
SocketMode=0600

[Install]
WantedBy=sockets.target
